# System/service/SD_touchdesk.service.d/override.conf
[Unit]
After=splashscreen.service
Wants=splashscreen.service

[Service]
# Очистить сплэш до старта приложения
ExecStartPre=/bin/bash /opt/splash/clear-splash.sh